# src/docker/config/stanis-base.service
[Unit]
Description=Stanis base service
Requires=docker.service
After=docker.service

[Service]
User=
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=
ExecStart=/usr/local/bin/docker-compose up -d stanis-base
ExecStop=/usr/local/bin/docker-compose stop stanis-base
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
